# systemd/.config/systemd/user/macastd.service
[Unit]
Description=macast Daemon

[Service]
Type=simple
ExecStart=/usr/bin/macast

[Install]
WantedBy=default.target
